# launchd/com.nanoclaw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanogemclaw</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/DevDisk/NanoGemClaw/launchd/run.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/DevDisk/NanoGemClaw</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/johnsontsai/.nvm/versions/node/v24.12.0/bin:/Users/johnsontsai/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/johnsontsai</string>
        <key>ASSISTANT_NAME</key>
        <string>Dart</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Volumes/DevDisk/NanoGemClaw/logs/nanoclaw.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/DevDisk/NanoGemClaw/logs/nanoclaw.error.log</string>
</dict>
</plist>
